# scripts/systemd/backup.service
[Unit]
Description=Backup Islandora

[Service]
Type=oneshot
ExecStart=/opt/islandora/d10_lehigh_agile/scripts/maintenance/backup.sh
