# scripts/systemd/isle-cron.timer
[Unit]
Description=ISLE cron timer

[Timer]
OnCalendar=*-*-* *:20:00

[Install]
WantedBy=timers.target
